# i386-squashfs-root/var/lib/dpkg/info/apport.list
/.
/etc
/etc/apport
/etc/apport/blacklist.d
/etc/apport/blacklist.d/README.blacklist
/etc/apport/blacklist.d/apport
/etc/default
/etc/default/apport
/etc/cron.daily
/etc/cron.daily/apport
/etc/bash_completion.d
/etc/bash_completion.d/apport_completion
/etc/init
/etc/init/apport.conf
/etc/init.d
/etc/logrotate.d
/etc/logrotate.d/apport
/usr
/usr/share
/usr/share/apport
/usr/share/apport/apport
/usr/share/apport/apport-checkreports
/usr/share/apport/package_hook
/usr/share/apport/kernel_crashdump
/usr/share/apport/kernel_oops
/usr/share/apport/gcc_ice_hook
/usr/share/apport/apportcheckresume
/usr/share/apport/unkillable_shutdown
/usr/share/apport/dump_acpi_tables.py
/usr/share/apport/java_uncaught_exception
/usr/share/apport/testsuite
/usr/share/apport/testsuite/crash
/usr/share/apport/testsuite/crash-digger
/usr/share/apport/testsuite/hooks
/usr/share/apport/testsuite/java
/usr/share/apport/testsuite/run
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_ubiquity.py
/usr/share/apport/package-hooks/source_linux.py
/usr/share/apport/package-hooks/source_apport.py
/usr/share/apport/package-hooks/source_debian-installer.py
/usr/share/apport/general-hooks
/usr/share/apport/general-hooks/generic.py
/usr/share/apport/general-hooks/automatix.py
/usr/share/apport/general-hooks/parse_segv.py
/usr/share/apport/general-hooks/ubuntu.py
/usr/share/apport/apport.jar
/usr/share/doc
/usr/share/doc/apport
/usr/share/doc/apport/crashdb-conf.txt
/usr/share/doc/apport/NEWS.gz
/usr/share/doc/apport/symptoms.txt
/usr/share/doc/apport/README
/usr/share/doc/apport/copyright
/usr/share/doc/apport/package-hooks.txt.gz
/usr/share/doc/apport/changelog.Debian.gz
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/scalable
/usr/share/icons/hicolor/scalable/mimetypes
/usr/share/icons/hicolor/scalable/apps
/usr/share/icons/hicolor/scalable/apps/apport.svg
/usr/share/mime
/usr/share/mime/packages
/usr/share/mime/packages/apport.xml
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/apport-cli.1.gz
/usr/share/man/man1/apport-unpack.1.gz
/usr/share/man/man1/apport-bug.1.gz
/usr/share/python
/usr/share/python/runtime.d
/usr/share/python/runtime.d/apport.rtupdate
/usr/bin
/usr/bin/apport-cli
/usr/bin/apport-unpack
/usr/bin/apport-bug
/etc/init.d/apport
/usr/share/apport/package-hooks/source_linux-meta.py
/usr/share/icons/hicolor/scalable/mimetypes/text-x-apport.svg
/usr/share/man/man1/apport-collect.1.gz
/usr/share/man/man1/ubuntu-bug.1.gz
/usr/bin/apport-collect
/usr/bin/ubuntu-bug
