# i386-squashfs-root/var/lib/dpkg/info/apport.conffiles
/etc/apport/blacklist.d/README.blacklist
/etc/apport/blacklist.d/apport
/etc/default/apport
/etc/cron.daily/apport
/etc/bash_completion.d/apport_completion
/etc/init/apport.conf
/etc/logrotate.d/apport
