# i386-squashfs-root/var/lib/dpkg/info/bluez.list
/.
/usr
/usr/bin
/usr/bin/bluez-test-device
/usr/bin/bluez-test-input
/usr/bin/bluez-test-network
/usr/bin/bluez-simple-service
/usr/bin/dfutool
/usr/bin/bluez-test-adapter
/usr/bin/bluez-simple-agent
/usr/bin/bluez-test-service
/usr/bin/bluez-test-serial
/usr/bin/sdptool
/usr/bin/bluez-test-telephony
/usr/bin/bluez-test-manager
/usr/bin/hcitool
/usr/bin/l2ping
/usr/bin/bluetooth-agent
/usr/bin/rctest
/usr/bin/bluez-test-audio
/usr/bin/ciptool
/usr/bin/bluez-test-discovery
/usr/bin/rfcomm
/usr/bin/l2test
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_bluez.py
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/bluez
/usr/share/doc
/usr/share/doc/bluez
/usr/share/doc/bluez/manager-api.txt
/usr/share/doc/bluez/examples
/usr/share/doc/bluez/examples/README
/usr/share/doc/bluez/examples/list-devices
/usr/share/doc/bluez/examples/monitor-bluetooth
/usr/share/doc/bluez/hfp-api.txt
/usr/share/doc/bluez/input-api.txt
/usr/share/doc/bluez/service-api.txt
/usr/share/doc/bluez/copyright
/usr/share/doc/bluez/audio-api.txt.gz
/usr/share/doc/bluez/device-api.txt.gz
/usr/share/doc/bluez/media-api.txt.gz
/usr/share/doc/bluez/control-api.txt
/usr/share/doc/bluez/network-api.txt
/usr/share/doc/bluez/sap-api.txt
/usr/share/doc/bluez/serial-api.txt
/usr/share/doc/bluez/attribute-api.txt.gz
/usr/share/doc/bluez/health-api.txt.gz
/usr/share/doc/bluez/adapter-api.txt.gz
/usr/share/doc/bluez/agent-api.txt
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/hciemu.1.gz
/usr/share/man/man1/bluetooth-agent.1.gz
/usr/share/man/man1/hcitool.1.gz
/usr/share/man/man1/sdptool.1.gz
/usr/share/man/man1/rfcomm.1.gz
/usr/share/man/man1/hcitrace.1.gz
/usr/share/man/man1/dfutool.1.gz
/usr/share/man/man1/ciptool.1.gz
/usr/share/man/man1/bluez-simple-agent.1.gz
/usr/share/man/man1/rctest.1.gz
/usr/share/man/man1/l2test.1.gz
/usr/share/man/man8
/usr/share/man/man8/bluetoothd.8.gz
/usr/share/man/man8/hid2hci.8.gz
/usr/share/man/man8/bccmd.8.gz
/usr/share/man/man8/l2ping.8.gz
/usr/share/man/man8/hciconfig.8.gz
/usr/share/man/man8/hciattach.8.gz
/usr/sbin
/usr/sbin/bccmd
/usr/sbin/hciemu
/usr/sbin/bluetoothd
/usr/sbin/hciattach
/usr/sbin/hciconfig
/usr/sbin/hcitrace
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/62-bluez-hid2hci.rules
/lib/udev/hid2hci
/lib/udev/bluez
/etc
/etc/default
/etc/default/bluetooth
/etc/init.d
/etc/init.d/bluetooth
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/bluetooth.conf
/etc/bluetooth
/etc/bluetooth/input.conf
/etc/bluetooth/network.conf
/etc/bluetooth/audio.conf
/etc/bluetooth/rfcomm.conf
/etc/bluetooth/serial.conf
/etc/bluetooth/main.conf
/usr/share/doc/bluez/README.Debian.gz
/usr/share/doc/bluez/README
/usr/share/doc/bluez/AUTHORS
/usr/share/doc/bluez/NEWS.Debian.gz
/usr/share/doc/bluez/TODO.gz
/usr/share/doc/bluez/changelog.Debian.gz
