# i386-squashfs-root/var/lib/dpkg/info/cryptsetup.conffiles
/etc/default/cryptdisks
/etc/init.d/cryptdisks-early
/etc/init.d/cryptdisks
/etc/init/cryptdisks-enable.conf
/etc/init/cryptdisks-udev.conf
