# i386-squashfs-root/var/lib/dpkg/info/friendly-recovery.list
/.
/lib
/lib/recovery-mode
/lib/recovery-mode/recovery-menu
/lib/recovery-mode/l10n.sh
/lib/recovery-mode/options
/lib/recovery-mode/options/grub
/lib/recovery-mode/options/root
/lib/recovery-mode/options/dpkg
/lib/recovery-mode/options/remount
/lib/recovery-mode/options/fsck
/lib/recovery-mode/options/clean
/lib/recovery-mode/options/netroot
/lib/recovery-mode/options/apt-snapshots
/lib/recovery-mode/options/failsafeX
/usr
/usr/share
/usr/share/doc
/usr/share/doc/friendly-recovery
/usr/share/doc/friendly-recovery/README
/usr/share/doc/friendly-recovery/TODO
/usr/share/doc/friendly-recovery/copyright
/usr/share/doc/friendly-recovery/changelog.gz
/etc
/etc/init
/etc/init/friendly-recovery.conf
/etc/init.d
/usr/share/recovery-mode
/etc/init.d/friendly-recovery
