# i386-squashfs-root/var/lib/dpkg/info/grub-common.list
/.
/usr
/usr/sbin
/usr/sbin/grub-mkconfig
/usr/sbin/grub-mkdevicemap
/usr/sbin/grub-probe
/usr/share
/usr/share/bug
/usr/share/bug/grub-common
/usr/share/bug/grub-common/presubj
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_grub2.py
/usr/share/grub
/usr/share/grub/unicode.pf2
/usr/share/grub/ascii.pf2
/usr/share/grub/euro.pf2
/usr/share/grub/ascii.h
/usr/share/grub/widthspec.h
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/grub-fstest.1.gz
/usr/share/man/man1/grub-menulst2cfg.1.gz
/usr/share/man/man1/grub-mkfont.1.gz
/usr/share/man/man1/grub-mkimage.1.gz
/usr/share/man/man1/grub-mklayout.1.gz
/usr/share/man/man1/grub-mkrelpath.1.gz
/usr/share/man/man1/grub-mount.1.gz
/usr/share/man/man1/grub-script-check.1.gz
/usr/share/man/man1/grub-mkrescue.1.gz
/usr/share/man/man1/grub-bin2h.1.gz
/usr/share/man/man1/grub-editenv.1.gz
/usr/share/man/man1/grub-mkpasswd-pbkdf2.1.gz
/usr/share/man/man8
/usr/share/man/man8/grub-probe.8.gz
/usr/share/man/man8/grub-mkconfig.8.gz
/usr/share/man/man8/grub-mkdevicemap.8.gz
/usr/share/doc
/usr/share/doc/grub-common
/usr/share/doc/grub-common/AUTHORS
/usr/share/doc/grub-common/README
/usr/share/doc/grub-common/THANKS
/usr/share/doc/grub-common/TODO
/usr/share/doc/grub-common/copyright
/usr/share/doc/grub-common/NEWS.Debian.gz
/usr/share/doc/grub-common/examples
/usr/share/doc/grub-common/examples/grub.cfg
/usr/share/doc/grub-common/NEWS.gz
/usr/share/doc/grub-common/changelog.Debian.gz
/usr/bin
/usr/bin/grub-bin2h
/usr/bin/grub-editenv
/usr/bin/grub-fstest
/usr/bin/grub-kbdcomp
/usr/bin/grub-menulst2cfg
/usr/bin/grub-mkfont
/usr/bin/grub-mkimage
/usr/bin/grub-mklayout
/usr/bin/grub-mkpasswd-pbkdf2
/usr/bin/grub-mkrelpath
/usr/bin/grub-mount
/usr/bin/grub-script-check
/usr/bin/grub-mkrescue
/usr/lib
/usr/lib/grub
/usr/lib/grub/grub-mkconfig_lib
/usr/lib/grub/update-grub_lib
/etc
/etc/pm
/etc/pm/sleep.d
/etc/pm/sleep.d/10_grub-common
/etc/grub.d
/etc/grub.d/05_debian_theme
/etc/grub.d/README
/etc/grub.d/00_header
/etc/grub.d/10_linux
/etc/grub.d/20_linux_xen
/etc/grub.d/30_os-prober
/etc/grub.d/40_custom
/etc/grub.d/41_custom
/etc/bash_completion.d
/etc/bash_completion.d/grub
/etc/init.d
/etc/init.d/grub-common
