# i386-squashfs-root/var/lib/dpkg/info/lightdm.list
/.
/var
/var/log
/var/log/lightdm
/var/cache
/var/cache/lightdm
/usr
/usr/sbin
/usr/sbin/lightdm
/usr/sbin/lightdm-session
/usr/sbin/guest-account
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/lightdm.1.gz
/usr/share/doc
/usr/share/doc/lightdm
/usr/share/doc/lightdm/lightdm.conf
/usr/share/doc/lightdm/keys.conf
/usr/share/doc/lightdm/copyright
/usr/share/doc/lightdm/changelog.Debian.gz
/usr/lib
/usr/lib/lightdm
/usr/lib/lightdm/lightdm-guest-session-wrapper
/usr/lib/lightdm/lightdm
/usr/lib/lightdm/lightdm/gdmflexiserver
/usr/lib/lightdm/lightdm-set-defaults
/usr/lib/lightdm/lightdm-greeter-session
/usr/bin
/usr/bin/dm-tool
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.DisplayManager.conf
/etc/lightdm
/etc/lightdm/users.conf
/etc/apparmor.d
/etc/apparmor.d/lightdm-guest-session
/etc/init
/etc/init/lightdm.conf
/etc/init.d
/etc/pam.d
/etc/pam.d/lightdm
/etc/pam.d/lightdm-autologin
/etc/init.d/lightdm
