# i386-squashfs-root/var/lib/dpkg/info/network-manager.list
/.
/var
/var/lib
/var/lib/NetworkManager
/etc
/etc/NetworkManager
/etc/NetworkManager/system-connections
/etc/NetworkManager/dispatcher.d
/etc/NetworkManager/dispatcher.d/01ifupdown
/etc/NetworkManager/NetworkManager.conf
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/NetworkManager.conf
/etc/dbus-1/system.d/nm-dhcp-client.conf
/etc/dbus-1/system.d/nm-dispatcher.conf
/etc/dbus-1/system.d/nm-avahi-autoipd.conf
/etc/init
/etc/init/network-manager.conf
/etc/init.d
/usr
/usr/sbin
/usr/sbin/NetworkManager
/usr/bin
/usr/bin/nm-tool
/usr/bin/nm-online
/usr/bin/nmcli
/usr/lib
/usr/lib/NetworkManager
/usr/lib/NetworkManager/nm-crash-logger
/usr/lib/NetworkManager/nm-dhcp-client.action
/usr/lib/NetworkManager/nm-dispatcher.action
/usr/lib/NetworkManager/nm-avahi-autoipd.action
/usr/lib/NetworkManager/libnm-settings-plugin-ifupdown.so
/usr/lib/NetworkManager/ifblacklist_migrate.sh
/usr/lib/pppd
/usr/lib/pppd/2.4.5
/usr/lib/pppd/2.4.5/nm-pppd-plugin.so
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/nmcli.1.gz
/usr/share/man/man1/nm-tool.1.gz
/usr/share/man/man1/nm-online.1.gz
/usr/share/man/man5
/usr/share/man/man5/NetworkManager.conf.5.gz
/usr/share/man/man5/nm-system-settings.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/NetworkManager.8.gz
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.nm_dispatcher.service
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.NetworkManager.policy
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_network-manager.py
/usr/share/doc
/usr/share/doc/network-manager
/usr/share/doc/network-manager/changelog.Debian.gz
/usr/share/doc/network-manager/AUTHORS
/usr/share/doc/network-manager/README.Debian
/usr/share/doc/network-manager/copyright
/usr/share/doc/network-manager/NEWS.Debian.gz
/usr/share/doc/network-manager/NEWS.gz
/usr/share/doc/network-manager/README.gz
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/77-nm-olpc-mesh.rules
/etc/init.d/network-manager
/usr/share/apport/package-hooks/source_network-manager-applet.py
