# i386-squashfs-root/var/lib/dpkg/info/popularity-contest.list
/.
/usr
/usr/share
/usr/share/popularity-contest
/usr/share/popularity-contest/popcon-upload
/usr/share/popularity-contest/default.conf
/usr/share/doc
/usr/share/doc/popularity-contest
/usr/share/doc/popularity-contest/FAQ.gz
/usr/share/doc/popularity-contest/examples
/usr/share/doc/popularity-contest/examples/popcon-submit.cgi
/usr/share/doc/popularity-contest/examples/README.examples
/usr/share/doc/popularity-contest/examples/clean-genpkglist
/usr/share/doc/popularity-contest/examples/prepop.pl
/usr/share/doc/popularity-contest/examples/popcon-process.sh
/usr/share/doc/popularity-contest/examples/popcon-stat.pl
/usr/share/doc/popularity-contest/examples/popcon.pl
/usr/share/doc/popularity-contest/examples/popanal.py
/usr/share/doc/popularity-contest/examples/clean-filter
/usr/share/doc/popularity-contest/copyright
/usr/share/doc/popularity-contest/README.gz
/usr/share/doc/popularity-contest/changelog.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/popularity-contest.8.gz
/usr/share/man/man8/popcon-largest-unused.8.gz
/usr/sbin
/usr/sbin/popularity-contest
/usr/sbin/popcon-largest-unused
/etc
/etc/cron.daily
/etc/cron.daily/popularity-contest
