# i386-squashfs-root/var/lib/dpkg/info/pppconfig.list
/.
/usr
/usr/share
/usr/share/menu
/usr/share/menu/pppconfig
/usr/share/doc
/usr/share/doc/pppconfig
/usr/share/doc/pppconfig/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/pppconfig.8.gz
/usr/share/man/fr
/usr/share/man/fr/man8
/usr/share/man/fr/man8/pppconfig.8.gz
/usr/sbin
/usr/sbin/pppconfig
/var
/var/cache
/var/cache/pppconfig
/etc
/etc/ppp
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/0dns-up
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/0dns-down
/etc/ppp/resolv
/etc/init.d
/etc/init.d/dns-clean
