# i386-squashfs-root/var/lib/dpkg/info/sane-utils.list
/.
/usr
/usr/bin
/usr/bin/umax_pp
/usr/bin/gamma4scanimage
/usr/bin/scanimage
/usr/bin/sane-find-scanner
/usr/share
/usr/share/doc
/usr/share/doc/sane-utils
/usr/share/doc/sane-utils/copyright
/usr/share/doc/sane-utils/README.Debian
/usr/share/doc/sane-utils/changelog.Debian.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/sane-find-scanner.1.gz
/usr/share/man/man1/gamma4scanimage.1.gz
/usr/share/man/man1/scanimage.1.gz
/usr/share/man/man8
/usr/share/man/man8/saned.8.gz
/usr/sbin
/usr/sbin/saned
/etc
/etc/default
/etc/default/saned
/etc/init.d
/etc/init.d/saned
/etc/sane.d
/etc/sane.d/saned.conf
