# i386-squashfs-root/etc/update-motd.d/91-release-upgrade
#!/bin/sh

if [ -x /usr/lib/update-manager/release-upgrade-motd ]; then
    exec /usr/lib/update-manager/release-upgrade-motd
fi
